# scripts/v2.0/deploy/vultr/defirisk-webhook.service
[Unit]
Description=DeFi Risk Webhook API (gunicorn)
After=network.target

[Service]
Type=simple
User=defirisk
Group=defirisk
WorkingDirectory=/opt/defi-risk/app/scripts/v2.0
EnvironmentFile=/etc/defi-risk/webhook.env
ExecStart=/opt/defi-risk/venv/bin/gunicorn --workers 2 --threads 4 --bind 127.0.0.1:5001 --timeout 180 webhook_server:app
Restart=always
RestartSec=3
StandardOutput=append:/var/log/defi-risk/webhook.log
StandardError=append:/var/log/defi-risk/webhook.err.log

[Install]
WantedBy=multi-user.target
