# scripts/v2.0/web_portal/deploy/security/cloudflare-ufw-sync.timer
[Unit]
Description=Run Cloudflare UFW sync monthly

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=2h
Unit=cloudflare-ufw-sync.service

[Install]
WantedBy=timers.target
