# scripts/v2.0/web_portal/deploy/security/crowdsec-slack-alerts.timer
[Unit]
Description=Run CrowdSec Slack digest every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Unit=crowdsec-slack-alerts.service
Persistent=true

[Install]
WantedBy=timers.target

