# scripts/v2.0/web_portal/deploy/security/crowdsec-slack-alerts.service
[Unit]
Description=Hodler Suite CrowdSec Slack Alert Digest
After=network-online.target crowdsec.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/crowdsec-slack-alerts
ExecStart=/usr/local/sbin/crowdsec_slack_alerts.sh

[Install]
WantedBy=multi-user.target

