# scripts/v2.0/web_portal/deploy/security/uptimerobot-slack-relay.service
[Unit]
Description=Hodler Suite UptimeRobot to Slack relay
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/uptimerobot-slack-relay
ExecStart=/usr/local/sbin/uptimerobot_slack_relay.sh
User=root
