# scripts/v2.0/web_portal/deploy/support/support-resend-sync.timer
[Unit]
Description=Run Hodler Suite Resend inbound fallback sync every 2 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
RandomizedDelaySec=20s
Persistent=true
Unit=support-resend-sync.service

[Install]
WantedBy=timers.target

