# scripts/v2.0/web_portal/deploy/privacy/user-telemetry-retention.service
[Unit]
Description=Hodler Suite user telemetry retention purge
After=network-online.target hodler-web-portal.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hodler_user_telemetry_retention.sh
User=root
Group=root
Nice=10
