# scripts/v2.0/web_portal/deploy/risk/risk-worker.timer
[Unit]
Description=Run Hodler Suite risk worker every 20 seconds

[Timer]
OnBootSec=35s
OnUnitActiveSec=20s
AccuracySec=2s
Unit=risk-worker.service
Persistent=true

[Install]
WantedBy=timers.target

