# scripts/v2.0/web_portal/deploy/risk/risk-worker.service
[Unit]
Description=Hodler Suite Risk worker (run one queued job)
After=network-online.target hodler-web-portal.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hodler_risk_worker_run_once.sh

